# src/main/resources/linux/crontab
## auto backup db
0 3 * * * /opt/bin/backup_db.sh